# raspberry-pi-2/rules.d/90-local.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0403", ATTR{idProduct}=="6001", RUN+="/usr/bin/sudo /home/pi/bin/start_tunslip.sh"
